# deploy/systemd/okr-kpi-frontend.service
[Unit]
Description=OKR KPI Frontend Preview (Optional)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/okr-kpi/frontend
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run preview -- --host 127.0.0.1 --port 4173
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
