# deploy/systemd/okr-kpi-frontend.service
[Unit]
Description=OKR KPI Frontend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/okr-kpi/frontend
Environment=NODE_ENV=production
EnvironmentFile=-/var/www/okr-kpi/frontend/.env.local
ExecStart=/usr/bin/npm run start -- --hostname 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
